# wizard/debian/prerm
#!/bin/bash

set -eo pipefail

case "$1" in
    remove|upgrade|deconfigure)
        # Stop services if they are running
        if systemctl is-active --quiet phoenixd-bebop 2>/dev/null; then
            echo "Stopping phoenixd service..."
            systemctl stop phoenixd-bebop || true
        fi

        if systemctl is-active --quiet minio-bebop 2>/dev/null; then
            echo "Stopping MinIO service..."
            systemctl stop minio-bebop || true
        fi

        if systemctl is-active --quiet bebop 2>/dev/null; then
            echo "Stopping be-BOP service..."
            systemctl stop bebop || true
        fi

        # Disable services if they are enabled
        if systemctl is-enabled --quiet phoenixd-bebop 2>/dev/null; then
            echo "Disabling phoenixd service..."
            systemctl disable phoenixd-bebop || true
        fi

        if systemctl is-enabled --quiet minio-bebop 2>/dev/null; then
            echo "Disabling MinIO service..."
            systemctl disable minio-bebop || true
        fi

        if systemctl is-enabled --quiet bebop 2>/dev/null; then
            echo "Disabling be-BOP service..."
            systemctl disable bebop || true
        fi
        ;;

    failed-upgrade)
        ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# Code autogenerated by debhelper goes here
#DEBHELPER#

exit 0
